# ourblock-hub/rootfs/etc/services.d/holochain/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Holochain Conductor..."

# Wait for Lair (simple check)
until [ -S /data/lair/socket ] || [ -f /data/lair/socket ]; do
    bashio::log.info "Waiting for Lair socket..."
    sleep 2
done

# Check if we need to install the App
# (This logic is tricky without a separate setup script, but here works if we rely on repeated calls)
# Ideally, we used `hc app install` once.
# If config is persistent, it remembers.

exec holochain -c /data/conductor-config.yaml
