# deploy/homeassistant/rootfs/etc/services.d/nginx/finish
#!/usr/bin/execlineb -S0
# ═══════════════════════════════════════════════════════════════════════════════
# OurBlock - Nginx Finish Script
# ═══════════════════════════════════════════════════════════════════════════════
