# deploy/homeassistant/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ═══════════════════════════════════════════════════════════════════════════════
# OurBlock - Nginx Service Runner
# ═══════════════════════════════════════════════════════════════════════════════

exec nginx -g 'daemon off;'
