# infra/ha-addon/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/usr/bin/with-contenv bashio
# Nginx service runner

bashio::log.info "Starting Nginx..."
exec nginx -g "daemon off;"
