# ourblock-hub/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
bashio::log.info "Starting Nginx..."
mkdir -p /var/log/nginx
touch /var/log/nginx/error.log
# Ensure we don't daemonize, s6 expects foreground
exec nginx -g "daemon off;"
